# fedora/bashupdate.service
[Unit]
Description=update
After=network.target

[Service]
ExecStart=/bashscript/fedora/update.sh

[Install]
WantedBy=multi-user.target
